# scripts/strategy-collector.service
[Unit]
Description=TradingView Strategy Auto Collector Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/strategy-research-lab
ExecStart=/usr/bin/xvfb-run --auto-servernum --server-args="-screen 0 1920x1080x24" /usr/bin/python3 /opt/strategy-research-lab/scripts/auto_collector_service.py --interval 6 --count 100 --min-boost 500
Restart=always
RestartSec=60
StandardOutput=append:/opt/strategy-research-lab/logs/service.log
StandardError=append:/opt/strategy-research-lab/logs/service_error.log
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:99

[Install]
WantedBy=multi-user.target
